# midas.service
[Unit]
Description=Midas — Autonomous BTC Trading Agent
After=network.target

[Service]
Type=simple
User=kyle
WorkingDirectory=/home/kyle/Midas
EnvironmentFile=/home/kyle/Midas/.env
ExecStart=/home/kyle/Midas/venv/bin/python main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
